# deploy/nayzfreedom-production-summary.service
[Unit]
Description=NayzFreedom daily production summary (one-shot)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=nayzfreedom
WorkingDirectory=/opt/nayzfreedom
EnvironmentFile=/opt/nayzfreedom/.env
ExecStart=/opt/nayzfreedom/.venv/bin/python /opt/nayzfreedom/production_summary.py
