# deploy/nayzfreedom-backup.timer
[Unit]
Description=NayzFreedom production backup — daily 03:30 UTC
Requires=nayzfreedom-backup.service

[Timer]
OnCalendar=*-*-* 03:30:00 UTC
AccuracySec=5min
Persistent=true

[Install]
WantedBy=timers.target
