# deploy/nayzfreedom-reporter.service
[Unit]
Description=NayzFreedom Weekly Performance Reporter (one-shot)
After=network.target

[Service]
Type=oneshot
User=nayzfreedom
WorkingDirectory=/opt/nayzfreedom
EnvironmentFile=/opt/nayzfreedom/.env
ExecStart=/opt/nayzfreedom/.venv/bin/python reporter.py
StandardOutput=journal
StandardError=journal
