# deploy/nayzfreedom-dashboard.service
[Unit]
Description=NayzFreedom Dashboard
After=network.target
Wants=network.target

[Service]
Type=simple
User=nayzfreedom
WorkingDirectory=/opt/nayzfreedom
EnvironmentFile=/opt/nayzfreedom/.env
ExecStart=/opt/nayzfreedom/.venv/bin/python dashboard.py --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
